# docker/prod/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /code/docker/prod/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
